# infrastructure/supervisor/atlas-worker.conf
; Atlas queue workers — Supervisor configuration
; Five named queues, each processed by a dedicated worker group.
; Adjust numprocs based on available CPU and expected throughput.
;
; IMPORTANT — queue:work's first positional argument is a CONNECTION name,
; not a queue name. Every job in this app dispatches via $this->onQueue(...)
; only (no job ever calls onConnection()), so every job is pushed onto
; whatever QUEUE_CONNECTION resolves to (the "database" connection by
; default — see .env.example) with its queue column set to high/ai/default/
; observations/maintenance. config/queue.php also defines four same-named
; connections using the redis driver ("Atlas queue topology" block) — those
; are a DIFFERENT, currently-unused routing path: nothing in this codebase
; calls onConnection('high') etc., so a worker started as `queue:work high`
; (passing "high" as the connection) listens on Redis's empty 'high' list
; and will never see a job that was actually pushed to the database
; connection's `jobs` table. Confirmed by direct dispatch-and-consume test
; (see docs/deployment/Queue-Workers.md, SCRUM-41) — the previous version of
; this file used exactly that broken form and would never have processed a
; single production job. Always pass the queue name via --queue=, with NO
; connection positional argument, so the worker uses config('queue.default')
; — the same connection every job actually dispatches through.

[program:atlas-worker-high]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/atlas/backend/artisan queue:work --queue=high --sleep=1 --tries=3 --max-time=3600
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=1
redirect_stderr=true
stdout_logfile=/var/log/atlas/worker-high.log
stopwaitsecs=60

[program:atlas-worker-ai]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/atlas/backend/artisan queue:work --queue=ai --sleep=3 --tries=3 --max-time=3600
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=2
redirect_stderr=true
stdout_logfile=/var/log/atlas/worker-ai.log
stopwaitsecs=300

[program:atlas-worker-default]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/atlas/backend/artisan queue:work --queue=default --sleep=3 --tries=3 --max-time=3600
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=2
redirect_stderr=true
stdout_logfile=/var/log/atlas/worker-default.log
stopwaitsecs=60

[program:atlas-worker-observations]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/atlas/backend/artisan queue:work --queue=observations --sleep=5 --tries=3 --max-time=3600
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=1
redirect_stderr=true
stdout_logfile=/var/log/atlas/worker-observations.log
stopwaitsecs=120

[program:atlas-worker-maintenance]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/atlas/backend/artisan queue:work --queue=maintenance --sleep=30 --tries=1 --max-time=3600
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=1
redirect_stderr=true
stdout_logfile=/var/log/atlas/worker-maintenance.log
stopwaitsecs=60
