# infrastructure/supervisor/atlas-worker.conf
; Atlas queue workers — Supervisor configuration
; Five named queues, each processed by a dedicated worker group.
; Adjust numprocs based on available CPU and expected throughput.

[program:atlas-worker-high]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/atlas/backend/artisan queue:work high --sleep=1 --tries=3 --max-time=3600
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=1
redirect_stderr=true
stdout_logfile=/var/log/atlas/worker-high.log
stopwaitsecs=60

[program:atlas-worker-ai]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/atlas/backend/artisan queue:work ai --sleep=3 --tries=3 --max-time=3600
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=2
redirect_stderr=true
stdout_logfile=/var/log/atlas/worker-ai.log
stopwaitsecs=300

[program:atlas-worker-default]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/atlas/backend/artisan queue:work default --sleep=3 --tries=3 --max-time=3600
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=2
redirect_stderr=true
stdout_logfile=/var/log/atlas/worker-default.log
stopwaitsecs=60

[program:atlas-worker-observations]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/atlas/backend/artisan queue:work observations --sleep=5 --tries=3 --max-time=3600
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=1
redirect_stderr=true
stdout_logfile=/var/log/atlas/worker-observations.log
stopwaitsecs=120

[program:atlas-worker-maintenance]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/atlas/backend/artisan queue:work maintenance --sleep=30 --tries=1 --max-time=3600
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=1
redirect_stderr=true
stdout_logfile=/var/log/atlas/worker-maintenance.log
stopwaitsecs=60
